# root/etc/s6-overlay/s6-rc.d/tfm-init/up
with-contenv

/etc/s6-overlay/s6-rc.d/tfm-init/tfm-init.sh
